# pm5.service
[Unit]
Description=PM5 Service
Requires=bluetooth.target network.target
After=bluetooth.target network.target

[Service]
User=root
Group=root
Type=simple
Environment="PM5_CLIENT_SECRET=abc123xyz"
ExecStart=/home/pi/go/bin/pm5 --logfile /var/log/pm5.log --db /var/spool/pm5/pm5.boltdb
ExecStop=/usr/bin/sleep 5
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
